# openvidu-server/docker/openvidu-server-kms-dev/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log;
pidfile=/var/run/supervisord.pid;
loglevel=debug

[program:kms]
command=/bin/bash /kms.sh
redirect_stderr=true

[program:openvidu-server]
command=/bin/bash -c "java -jar -Dspring.config.additional-location=classpath:/application-container.properties -agentlib:jdwp=transport=dt_socket,server=y,suspend=n,address=10002 /openvidu-server.jar"
redirect_stderr=true
